# croneditor.service
[Unit]
Description=CronEditor Web Application
After=network.target

[Service]
Type=simple
User={{USER}}
Group={{USER}}
WorkingDirectory={{INSTALL_DIR}}
Environment=PATH={{INSTALL_DIR}}/.venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/{{USER}}
Environment=USER={{USER}}
ExecStart={{INSTALL_DIR}}/.venv/bin/python app.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target